# backend/supervisor.conf
[program:tianting]
command=uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
directory=/app
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/tianting/uvicorn.log
stderr_logfile=/var/log/tianting/uvicorn_error.log
environment=PYTHONPATH="/app"

[program:tianting-worker]
command=arq app.tasks.worker:WorkerSettings
directory=/app
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/tianting/worker.log
stderr_logfile=/var/log/tianting/worker_error.log
environment=PYTHONPATH="/app"
